# deploy/systemd/hugo-mcp-go.service
[Unit]
Description=hugo-mcp-go backend service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hugo-mcp
Group=hugo-mcp
WorkingDirectory=/var/lib/hugo-mcp-go
EnvironmentFile=-/etc/hugo-mcp-go/hugo-mcp-go.env
ExecStart=/srv/hugo-mcp-go/bin/hugo-mcp-go
Restart=on-failure
RestartSec=2
NoNewPrivileges=yes
StateDirectory=hugo-mcp-go
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
CapabilityBoundingSet=
AmbientCapabilities=
LockPersonality=yes
MemoryMax=512M
CPUQuota=50%
UMask=0027
StandardOutput=journal
StandardError=journal
ReadWritePaths=/var/lib/hugo-mcp-go
ReadOnlyPaths=/etc/hugo-mcp-go/secrets
IPAddressDeny=any
RestrictAddressFamilies=AF_UNIX

[Install]
WantedBy=multi-user.target
