# deploy/systemd/hugo-mcp-go-http.service
[Unit]
Description=hugo-mcp-go HTTP native backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hugo-mcp
Group=hugo-mcp
SupplementaryGroups=jm
WorkingDirectory=/home/jm/hugo-site
EnvironmentFile=/etc/hugo-mcp-go/hugo-mcp-go-http.env
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/srv/hugo-mcp-go/bin/hugo-mcp-go
Restart=on-failure
RestartSec=2
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/jm/hugo-site/content /home/jm/hugo-site/static /home/jm/hugo-site/public /home/jm/hugo-site/work /home/jm/hugo-site/resources /home/jm/hugo-site/.hugo_build.lock
CapabilityBoundingSet=
LockPersonality=yes
MemoryMax=512M
CPUQuota=50%
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
