# deploy/systemd/hugo-mcp-go-http.service
[Unit]
Description=hugo-mcp-go HTTP native backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hugo-mcp
Group=hugo-mcp
SupplementaryGroups=hugo-mcp-shim
WorkingDirectory=/var/lib/hugo-mcp-go
EnvironmentFile=/etc/hugo-mcp-go/hugo-mcp-go-http.env
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/srv/hugo-mcp-go/bin/hugo-mcp-go
Restart=on-failure
RestartSec=2
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/hugo-mcp-go
CapabilityBoundingSet=
LockPersonality=yes
MemoryMax=512M
CPUQuota=50%
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
